# data_backup.service
[Unit]
Description=Make Data-Backup

[Service]
ExecStart=/home/USER/backup_data.py
