# data_backup.service
[Unit]
Description=Make Data-Backup

[Service]
ExecStart=/home/Dennis1/backup_data.py
